# deploy/systemd/sslctl.service
[Unit]
Description=SSL Certificate Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/sslctl daemon
Restart=always
RestartSec=30
User=root
Group=root

# 工作目录
WorkingDirectory=/opt/sslctl

# 日志
StandardOutput=journal
StandardError=journal

# 安全限制
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/sslctl
ReadWritePaths=-/etc/nginx
ReadWritePaths=-/etc/apache2
ReadWritePaths=-/etc/httpd
ReadWritePaths=-/etc/letsencrypt

[Install]
WantedBy=multi-user.target
